# systemd/hopeturtle-oled-boot.service
[Unit]
Description=HopeTurtle OLED boot messages
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
User=hopeturtle
Group=hopeturtle
# 1) Waking
ExecStart=/usr/bin/python3 /home/hopeturtle/hopeturtle/src/oled_status.py boot-waking
# 2) Alive
ExecStart=/usr/bin/python3 /home/hopeturtle/hopeturtle/src/oled_status.py boot-alive
# 3) Show last known distance to Al Mawasi
ExecStart=/usr/bin/python3 /home/hopeturtle/hopeturtle/src/oled_status.py distance
# 4) Optional: celebratory swim
# ExecStart=/usr/bin/python3 /home/hopeturtle/hopeturtle/src/oled_status.py swim

[Install]
WantedBy=multi-user.target
