# systemd/hopeturtle-gps.timer
# /etc/systemd/system/hopeturtle-gps.timer
#      ___________    _____
#  /               \ |  o  | 
# |                |/   __\| 
# |  _______________  /     
#   |_|_|     |_|_|
#
# HopeTurtle Timer 🐢
# Triggers the hopeturtle-gps.service regularly.
# Runs once ~30 seconds after boot, then every 5 minutes (00, 05, 10…).

[Unit]
Description=Run HopeTurtle GPS logger every 5 minutes
After=network.target pigpiod.service
Wants=pigpiod.service

[Timer]
# Start ~30s after boot, repeat every 5 minutes
OnBootSec=30s
OnCalendar=*:0/5
AccuracySec=15s
Persistent=true
Unit=hopeturtle-gps.service

[Install]
WantedBy=timers.target
