# systemd/hopeturtle-gps.service
# /etc/systemd/system/hopeturtle-gps.service
#   _________     ____
#  /          \  |  0 |
# |            |/ ___\
# |____________|_/
#   |__|  |__|
#
# HopeTurtle GPS Service
# Runs gps_snapshot.py once (systemd timer calls this every 5 minutes).
# Reads GPS from soft-serial GPIO17 and logs CSV + JSON output.
# After each run, displays fix info briefly on the OLED.

[Unit]
Description=HopeTurtle GPS one-shot logger
After=network.target pigpiod.service
Wants=pigpiod.service

[Service]
Type=oneshot

# ---- Environment configuration ----
Environment=HT_GPS_MODE=soft
Environment=HT_GPS_SOFT_RX=17
Environment=HT_BAUD=9600
Environment=HT_DATA_DIR=/home/hopeturtle/hopeturtle/data
Environment=PYTHONUNBUFFERED=1

# ---- Working directory ----
WorkingDirectory=/home/hopeturtle/hopeturtle/src

# ---- Main command ----
ExecStart=/usr/bin/python3 /home/hopeturtle/hopeturtle/src/gps_snapshot.py

# ---- OLED feedback (show brief summary after logging) ----
ExecStartPost=/usr/bin/python3 /home/hopeturtle/hopeturtle/src/oled_status.py brief

# ---- Run as user ----
User=hopeturtle
Group=hopeturtle

# ---- Lower priority ----
Nice=10

# ---- Restart policy ----
# (only runs when timer calls it, so no persistent process)
RemainAfterExit=no
